# clash/clash.service
[Unit]
Description=Clash Daemon

[Service]
ExecStart=/usr/bin/clash -d /etc/clash/
Restart=on-failure

[Install]
WantedBy=multi-user.target
